# debian/cron.d
01 22 * * * postgres /usr/bin/pgdumper.pl >> /var/log/pgdumper/pgdumper.log 2>&1
